# debian/postinst
#!/bin/sh
set -e

# Enable user service for all users
deb-systemd-helper --user enable pip-os.service || true

exit 0